# deploy/systemd/lv-liquidator.service
[Unit]
Description=Liquidity Vacuum — Liquidator
After=network-online.target redis-server.service
Wants=network-online.target redis-server.service

[Service]
Type=simple
User=lv
Group=lv
WorkingDirectory=/opt/liquidity
EnvironmentFile=-/opt/liquidity/.env
ExecStart=/opt/liquidity/bin/liquidator
Restart=always
RestartSec=2
StartLimitIntervalSec=0
WatchdogSec=60
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
